# azure-stock-monitor/com.nick.marukyu-stock.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.nick.marukyu-stock</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/nick/Documents/Codex/2026-06-26/Agent-Skills/azure-stock-monitor/run-local-hourly.sh</string>
  </array>
  <key>StartInterval</key>
  <integer>3600</integer>
  <key>RunAtLoad</key>
  <true/>
  <key>StandardOutPath</key>
  <string>/Users/nick/Documents/Codex/2026-06-26/Agent-Skills/azure-stock-monitor/logs/launchd.out.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/nick/Documents/Codex/2026-06-26/Agent-Skills/azure-stock-monitor/logs/launchd.err.log</string>
</dict>
</plist>
